# rootfs/lite/etc/systemd/system/agentd.service
[Unit]
Description=Kiki Agent Daemon (Lite/ARM)
Documentation=https://docs.kiki-os.com/agentd
After=network.target local-fs.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/agentd --config /etc/kiki/agentd.toml
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Conservative memory limit for edge/embedded devices
MemoryMax=512M
CPUQuota=80%

StateDirectory=kiki
RuntimeDirectory=kiki
RuntimeDirectoryPreserve=yes

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/kiki /run/kiki

[Install]
WantedBy=multi-user.target
