# rootfs/server/etc/systemd/system/agentd.service.d/10-server.conf
# Server-role drop-in for agentd.
#
# Headless server nodes phone home to the cloud control plane. The fleet client
# is built into agentd; this drop-in is where VPS provisioning (cloud-init
# user-data) injects the enrollment credential and, if needed, overrides.
#
# Provisioning typically writes the fleet endpoint + token here, e.g.:
#   [Service]
#   Environment=KIKI_FLEET_TOKEN=<bearer-from-enrollment>
#
# The fleet endpoint itself lives in /etc/kiki/agentd.toml ([fleet] cloud_url),
# also set by provisioning. Left empty here so the base unit runs unchanged
# until a node is provisioned.
[Service]
Environment=KIKI_LOG_JSON=1
