# rootfs/desktop/etc/systemd/system/kiki-ded.service
[Unit]
Description=Kiki Desktop Environment
Documentation=https://docs.kiki-os.com/de
After=systemd-logind.service agentd.service
Requires=agentd.service

[Service]
Type=simple
# The DE daemon: runs the shell logic and brings up the Wayland compositor when
# a display is present.
ExecStart=/usr/bin/kiki-ded --socket /run/kiki/wayland.sock
Restart=on-failure
RestartSec=3

# Wayland runtime requirements
Environment=XDG_RUNTIME_DIR=/run/user/0
Environment=WAYLAND_DISPLAY=wayland-0

RuntimeDirectory=user/0
RuntimeDirectoryPreserve=yes

# GPU access
SupplementaryGroups=video render input

[Install]
WantedBy=graphical.target
